# deploy/vllm-laguna-watchdog.timer
[Unit]
Description=Run the vLLM Laguna inference watchdog every 5 minutes

[Timer]
OnBootSec=10min
OnUnitActiveSec=5min

[Install]
WantedBy=timers.target
